# images/rpmbuild/base/files/etc/profile.d/00-aliases.sh
alias rpmbuild='rpmbuild --define "_topdir /opt/pkgbox/rpmbuild"'
